# systems/server-desktop/proxmox/90-uinput-lxc.rules
# Allow the server-desktop LXC container (GID 100989 = container's uinput group) to
# open /dev/uinput so sunshine can create virtual input devices.
KERNEL=="uinput", SUBSYSTEM=="misc", GROUP="100989", MODE="0660"

# Set permissions on virtual input devices sunshine creates, so the container's muser
# (host uid 100000) and input group (host GID 100989) can read them for libinput.
# ATTRS{phys}=="" matches only virtual devices (no physical bus path).
SUBSYSTEM=="input", KERNEL=="event[0-9]*", ATTRS{phys}=="", OWNER="100000", GROUP="100989", MODE="0660"

# When sunshine's virtual keyboard appears, trigger a udev input rescan inside the
# container so libinput discovers the new devices. The container's udevd never receives
# the original kernel uevent (it fires in the host namespace), but nsenter lets us run
# udevadm trigger in the container's network namespace with host root privileges, which
# can write to sysfs and causes the container's udevd to notify libinput.
# $$ is udev's escape for a literal $, producing $(lxc-info ...) for the shell.
SUBSYSTEM=="input", ACTION=="add", ATTRS{name}=="Keyboard passthrough", RUN+="/bin/sh -c 'nsenter --target $$(lxc-info -H -p 117 2>/dev/null) --net -- udevadm trigger --action=add --subsystem-match=input 2>/dev/null'"
